# systemd/claude-config-sync.timer
[Unit]
Description=Claude Code config sync timer (safety net every 4 hours)

[Timer]
OnBootSec=5min
OnUnitActiveSec=4h
Persistent=true

[Install]
WantedBy=timers.target
